# docker/configs/supervisor/conf.d/spiderweb.conf
[program:spiderweb]
command=/usr/bin/sw-node || sleep 900
user=root
autostart=true
startretries=32
