# configs/supervisor/conf.d/api.conf
[program:api]
command=/node/sw-node || sleep 900
user=root
autostart=true
startretries=32
