# configs/supervisor/conf.d/rtunnel.conf
[program:rtunnel]
command=ssh -o ServerAliveInterval=300 -NR 2222:localhost:22 user@destination || sleep 900
user=root
autostart=false
startretries=32
