# worker_node/supervisor/kubelet.conf
[program:kubelet]
command=/usr/bin/kubelet --fail-swap-on=false
autostart=true
autorestart=true
