# worker_node/supervisor/tunnel.conf
[program:tunnel]
command=ssh -R 12345:localhost:22 ubuntu@spiderweb.com
autostart=true
autorestart=true
